# root/jdk-environ/etc/profile.d/jdk-environ.sh
export JAVA_HOME="/usr/lib/jvm/default"
export _JAVA_OPTIONS='-Dawt.useSystemAAFontSettings=on -Dswing.aatext=true'
